# crontabs/root
* * * * * /etc/init.d/lightning-rod restart

