# etc/cron.d/root_raspberry
* * * * * root /etc/init.d/s4t-lightning-rod start

