# etc/cron.d/root_yun
* * * * * /etc/init.d/lightning-rod start