# etc/cron.d/root_openwrt
* * * * * /etc/init.d/lightning-rod start
*/10 * * * *  ntpdate 0.pool.ntp.org